# config/systemd/openai-proxy.service
[Unit]
Description=OpenAI API Proxy
After=network.target
Documentation=https://github.com/ailocal/openai-proxy

[Service]
Type=forking
User=openai-proxy
Group=openai-proxy
RuntimeDirectory=openai-proxy
RuntimeDirectoryMode=0750

# Environment file location
EnvironmentFile=/etc/openai-proxy/env

# Paths
ExecStart=/usr/local/bin/openai-proxy start
ExecStop=/usr/local/bin/openai-proxy stop
ExecReload=/usr/local/bin/openai-proxy reload

# Give a reasonable amount of time for the server to start up/shut down
TimeoutSec=300

# Process management
Restart=on-failure
RestartSec=5s

# Hardening
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
